# logmeters/logmeters.socket
[Unit]
Description=pipe for logmeters

[Socket]
ListenFIFO=/run/user/1001/logmeters.socket

[Install]
WantedBy=default.target
