# ssd-health-dashboard/collector/store-smart.timer
[Unit]
Description=SMART data collector

[Timer]
OnBootSec=12h
OnUnitActiveSec=1d
AccuracySec=12h

[Install]
WantedBy=multi-user.target

